# cookbooks/edemocracia/templates/centos/celeryd.service.erb
[Unit]
Description=Celery worker
PartOf=colab.service

[Service]
EnvironmentFile=-/etc/default/celeryd
User=<%= node['config']['system']['user'] %>
ExecStart=<%= node['config']['colab']['virtualenv'] %>/bin/colab-admin celery worker $CELERYD_OPTS
Restart=on-failure

[Install]
WantedBy=multi-user.target