# cookbooks/discourse/templates/default/discourse-sidekiq.service.erb
[Unit]
Description=Discourse Sidekiq
After=syslog.target
After=network.target

[Service]
User=<%= node['config']['system']['user'] %>
Group=discourse
WorkingDirectory=<%= node['config']['discourse']['dir'] %>
ExecStart=/usr/local/bin/bundle exec sidekiq -e production
PIDFile=<%= node['config']['discourse']['dir'] %>/log/sk.pid

[Install]
WantedBy=multi-user.target