# cookbooks/wikilegis/templates/default/wikilegis.service.erb
[Unit]
Description=Wikilegis

[Service]
User=<%= node['config']['system']['user'] %>
EnvironmentFile=-/etc/sysconfig/wikilegis
PermissionsStartOnly=true
WorkingDirectory=<%= node['config']['wikilegis']['dir'] %>
ExecStartPre=/usr/bin/install -d -m 0750 -o <%= node['config']['system']['user'] %> -g wikilegis /var/lock/wikilegis
ExecStart=<%= node['config']['wikilegis']['virtualenv'] %>/bin/gunicorn wikilegis.wsgi:application --config=<%= node['config']['wikilegis']['dir'] %>/gunicorn.py
Restart=on-failure

[Install]
WantedBy=multi-user.target